# examples/systemd/mcp-native-guard-smoke.timer
[Unit]
Description=Schedule daily mcp-native-guard smoke test

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=15m
Unit=mcp-native-guard-smoke.service

[Install]
WantedBy=timers.target
